# php-fpm73/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=docker-php-entrypoint php-fpm
redirect_stderr=true

[program:cron]
command=/bin/bash -c "cron -f"
redirect_stderr=true