# php-fpm73/supervisord.conf
[program:php-cron]
command=/bin/bash -c "cron -f && php-fpm"
redirect_stderr=true